# puppy.service
[Unit]
Description=Puppy
After=network.target
[Service]
Type=simple
ExecStart=/puppy/puppy.sh
[Install]
WantedBy=network.target
